# .container/dev/podman-compose-dev.service
[Unit]
Description=Podman Compose Development Services
After=network.target
Requires=podman.socket

[Service]
Type=simple
# Update the path below to match your installation
WorkingDirectory=/path/to/monorepo/.container/dev
ExecStart=/usr/bin/podman compose up
ExecStop=/usr/bin/podman compose down
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300
TimeoutStopSec=60

# Run as your user (update USER and GROUP)
User=your-username
Group=your-groupname

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target


